# supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/django_tech_test/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
